# scripts/systemd/leninbot-stale-secrets.timer
[Unit]
Description=Leninbot stale-secrets weekly check

[Timer]
OnCalendar=Mon 10:00
Persistent=true

[Install]
WantedBy=timers.target
