# scripts/systemd/leninbot-stale-secrets.service
[Unit]
Description=Leninbot stale-secrets weekly check (notifies Telegram)
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/grass/leninbot
EnvironmentFile=/home/grass/leninbot/.env
LoadCredentialEncrypted=telegram_bot_token:/etc/credstore.encrypted/telegram_bot_token.cred
ExecStart=/home/grass/leninbot/venv/bin/python /home/grass/leninbot/scripts/manage_secrets.py list --stale 90 --notify
