# deploy/systemd/research-document-translation.timer
[Unit]
Description=Daily untranslated research document English translation

[Timer]
OnCalendar=*-*-* 04:20:00
Persistent=true
RandomizedDelaySec=20m
Unit=research-document-translation.service

[Install]
WantedBy=timers.target
