# systemd/leninbot-razvedchik.service
[Unit]
Description=Razvedchik — Moltbook patrol agent
After=network.target

[Service]
Type=oneshot
User=grass
WorkingDirectory=/home/grass/leninbot
ExecStart=/home/grass/leninbot/venv/bin/python agents/razvedchik/razvedchik.py --patrol --max-comments 3
EnvironmentFile=/home/grass/leninbot/.env
StandardOutput=journal
StandardError=journal
