# systemd/leninbot-telegram.service
[Unit]
Description=LeninBot Telegram Bot (aiogram)
After=network.target leninbot-neo4j.service leninbot-embedding.service
Wants=leninbot-neo4j.service leninbot-embedding.service

[Service]
User=grass
WorkingDirectory=/home/grass/leninbot
EnvironmentFile=/home/grass/leninbot/.env
ExecStart=/home/grass/leninbot/venv/bin/python -m telegram.bot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
