# systemd/leninbot-dream.timer
[Unit]
Description=Run Dream Diary daily at 03:00 KST (18:00 UTC)

[Timer]
OnCalendar=*-*-* 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
