# systemd/leninbot-autonomous.timer
[Unit]
Description=LeninBot Autonomous Project Timer (hourly, T0 pilot)

[Timer]
# Fires at :17 past every hour to avoid the usual background timer window.
# Adjust freely — each fire is one tick.
OnCalendar=*-*-* *:17:00 Asia/Seoul
Persistent=true

[Install]
WantedBy=timers.target
