# systemd/leninbot-commulingo-terms.timer
[Unit]
Description=Register one missing CommuLingo glossary term every fifteen minutes

[Timer]
OnCalendar=*:7/15
Persistent=true

[Install]
WantedBy=timers.target
