# systemd/leninbot-kg-integrity.service
[Unit]
Description=LeninBot KG integrity and search smoke check
After=network.target leninbot-neo4j.service
Requires=leninbot-neo4j.service

[Service]
Type=oneshot
User=grass
WorkingDirectory=/home/grass/leninbot
EnvironmentFile=/home/grass/leninbot/.env
LoadCredentialEncrypted=gemini_api_key:/etc/credstore.encrypted/gemini_api_key.cred
LoadCredentialEncrypted=neo4j_password:/etc/credstore.encrypted/neo4j_password.cred
LoadCredentialEncrypted=openai_api_key:/etc/credstore.encrypted/openai_api_key.cred
LoadCredentialEncrypted=telegram_bot_token:/etc/credstore.encrypted/telegram_bot_token.cred
ExecStart=/home/grass/leninbot/venv/bin/python scripts/check_kg_integrity.py --smoke-query DiaMat --notify
TimeoutStartSec=120
