# systemd/leninbot-kg-integrity.timer
[Unit]
Description=LeninBot KG integrity and search smoke check timer

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
