# systemd/leninbot-llama.service
[Unit]
Description=llama-server (qwen3.5-4b Q4_K_M) — leninbot 로컬 LLM 폴백
After=network.target

[Service]
Type=simple
User=grass
Environment=LD_LIBRARY_PATH=/home/grass/llama-cpp/llama-b8502
ExecStart=/home/grass/llama-cpp/llama-b8502/llama-server \
    -m /home/grass/llama-cpp/models/qwen3.5-4b-q4_k_m.gguf \
    --host 127.0.0.1 --port 11435 \
    -c 8192 -t 4
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
