# systemd/leninbot-variant-scan.service
[Unit]
Description=LeninBot CommuLingo name-variant scan (notifies Telegram on new candidates)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=grass
WorkingDirectory=/home/grass/leninbot
EnvironmentFile=/home/grass/leninbot/.env
Environment=PYTHONUNBUFFERED=1
LoadCredentialEncrypted=db_password:/etc/credstore.encrypted/db_password.cred
LoadCredentialEncrypted=telegram_bot_token:/etc/credstore.encrypted/telegram_bot_token.cred
ExecStart=/home/grass/leninbot/venv/bin/python /home/grass/leninbot/scripts/commulingo_find_name_variants.py --min-count 1 --notify
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal
