# systemd/leninbot-debate.service
[Unit]
Description=Cyber-Lenin AI Debate Generator
After=network.target

[Service]
Type=oneshot
User=grass
WorkingDirectory=/home/grass/leninbot
ExecStart=/home/grass/leninbot/venv/bin/python ai_debate.py
EnvironmentFile=/home/grass/leninbot/.env
StandardOutput=journal
StandardError=journal
