# systemd/leninbot-experience.timer
[Unit]
Description=LeninBot Experience Writer Timer (daily 00:30 KST)

[Timer]
OnCalendar=*-*-* 00:30:00 Asia/Seoul
Persistent=true

[Install]
WantedBy=timers.target
