# systemd/leninbot-riddle.timer
[Unit]
Description=Run Riddle Engine daily at 04:30 KST (19:30 UTC)

[Timer]
OnCalendar=*-*-* 19:30:00
Persistent=true

[Install]
WantedBy=timers.target
